# cronwowp2/crontab
01,31 * * * * python /cronwowp2/source/cron.py >> /var/log/cronwowp2.log

